# srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command = sh -c "vsftpd /etc/vsftpd/vsftpd.conf && kill -s SIGTERM $(cat supervisord.pid)"